# configs/archlenovop52/etc/systemd/system/update-system-config.service
[Unit]
Description=Update system configuration files
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/update-system-config
User=root
Group=root

[Install]
WantedBy=multi-user.target
