# configs/archlenovop52/etc/systemd/system/update-system-config.timer
[Unit]
Description=Run update-system-config.service daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
