# configs/archlenovop52/home/labin/.config/systemd/user/update-user-config.service
[Unit]
Description=Update user configuration files
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/update-user-config

[Install]
WantedBy=default.target
