# configs/archlenovop52/home/labin/.config/systemd/user/update-user-config.timer
[Unit]
Description=Run update-user-config.service daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
